# services/visual-break.service
[Unit]
Description=Visual break reminder

[Service]
Type=oneshot
ExecStart=/home/user/.local/bin/visual-break
